# 80-dvorak.rules
SUBSYSTEM=="input", ACTION=="add", RUN+="/bin/systemctl start dvorak@%k.service"
