# 80-dvorak.rules
SUBSYSTEM=="input", ACTION=="add", RUN+="/usr/bin/dvorak /dev/input/%k keyb &"
