# touchcursor.service
[Unit]
Description=Touch Cursor Service

[Service]
ExecStart=/usr/bin/touchcursor

[Install]
WantedBy=default.target