# touchcursor.service
[Unit]
Description=Touch Cursor Service
StartLimitIntervalSec=0
PartOf=touchcursor.path

[Service]
ExecStart=/usr/bin/touchcursor
