# touchcursor.path
[Unit]
Description=Touchcursor Path Monitor
StartLimitIntervalSec=0

[Path]
PathChanged=/dev/input

[Install]
WantedBy=default.target
