# touchcursor@.service
[Unit]
Description=Touch Cursor Service

[Service]
ExecStart=/usr/bin/touchcursor
User=%i
Group=

[Install]
WantedBy=multi-user.target