# touchcursor.service
[Unit]
Description=Touch Cursor Service

[Service]
ExecStart=%h/.local/bin/touchcursor
