# touchcursor-watcher.path
[Unit]
Description=Touchcursor Path Monitor
StartLimitIntervalSec=0

[Path]
Unit=touchcursor-watcher.service
PathChanged=/dev/input

[Install]
WantedBy=default.target
