# touchcursor.service
[Unit]
Description=Touch Cursor Service
StartLimitIntervalSec=0

[Service]
ExecStart=/usr/bin/touchcursor

[Install]
WantedBy=default.target