# touchcursor.path
[Unit]
Description=Touchcursor Path Monitor
StartLimitBurst=0
PropagatesStopTo=touchcursor.service

[Path]
Unit=touchcursor.service
PathChanged=/dev/input
PathModified=%h/.config/touchcursor/touchcursor.conf

[Install]
Also=touchcursor.service
WantedBy=default.target
